# data/thermald_msrfix.service.in
[Unit]
Description=Thermal Daemon MSR Fix
Before=thermald.service

[Service]
ExecStart=@bindir@/thermald_msrfix.sh

[Install]
WantedBy=thermald.service